# newrelic-plugin.conf
[program:newrelic-plugin]
user=root
command=newrelic-plugin-agent -f -c /etc/newrelic/newrelic-plugin-agent.cfg
autostart=false
priority=3
